# resources/ifmon.cron
#
# save network usage at every minute
#
*/1 * * * * root python /usr/share/ifmon/ifmon/ifmon.pyc

